# automation/system-services/linux/sophia-intel-ai.service
[Unit]
Description=Sophia Intel AI Platform
Documentation=https://github.com/sophia-intel-ai/docs
After=network.target docker.service postgresql.service redis.service
Wants=docker.service
RequiredBy=multi-user.target

[Service]
Type=exec
User=sophia
Group=sophia
WorkingDirectory=/opt/sophia-intel-ai
ExecStart=/opt/sophia-intel-ai/automation/scripts/startup-manager.sh --mode=development --environment=local --auto-recovery=true
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/opt/sophia-intel-ai/automation/scripts/startup-manager.sh --action=stop
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=300
TimeoutStopSec=60
Restart=on-failure
RestartSec=30
StartLimitBurst=3
StartLimitIntervalSec=300

# Environment
Environment=SOPHIA_ENVIRONMENT=development
Environment=SOPHIA_AUTO_RECOVERY=true
Environment=SOPHIA_LOG_LEVEL=INFO
Environment=SOPHIA_HOME=/opt/sophia-intel-ai
Environment=PYTHONPATH=/opt/sophia-intel-ai
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sophia-intel-ai/logs /opt/sophia-intel-ai/data /var/run/docker.sock
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=8192
LimitNPROC=4096
MemoryMax=4G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sophia-intel-ai

# Health checks
WatchdogSec=60

[Install]
WantedBy=multi-user.target
Alias=sophia.service
